# arch/supervisord.conf
[supervisord]
nodaemon=true

; [program:trace_streamer]
; command=sh -c "tail -F /var/log/envoy.log | python stream_traces.py"
; autostart=true
; autorestart=false
; startretries=3
; priority=1
; stdout_logfile=/dev/stdout
; stderr_logfile=/dev/stderr
; stdout_logfile_maxbytes = 0
; stderr_logfile_maxbytes = 0


[program:envoy]
command=sh -c "python config_generator.py && envsubst < /etc/envoy/envoy.yaml > /etc/envoy.env_sub.yaml && envoy -c /etc/envoy.env_sub.yaml --component-log-level wasm:debug 2>&1 | tee /var/log/envoy.log"
; command=sh -c "python config_generator.py && envsubst < /etc/envoy/envoy.yaml > /etc/envoy.env_sub.yaml && envoy -c /etc/envoy.env_sub.yaml --log-level trace 2>&1 | tee /var/log/envoy.log"
autostart=true
autorestart=true
startretries=3
priority=2
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
